# shellwords-api.service
[Unit]
Description=ShellWords API Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/html/shellwords/server
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=shellwords-api
Environment=NODE_ENV=production
Environment=PORT=3001

[Install]
WantedBy=multi-user.target
